# crontab
# Daily pipeline: fetch raw snapshots, then upsert them into Postgres.
# Times are UTC (container default).
30 5 * * * cd /app && .venv/bin/python -m pipeline collect && .venv/bin/python -m pipeline load
